# debian/nvidia-glx-dev.postinst
#!/bin/sh 

set -e


case "$1" in
	configure)
		if [ -e /etc/init.d/nvidia-glx ] ; then
			invoke-rc.d nvidia-glx start
		fi
		;;
esac		
	
#DEBHELPER#
